# modules/caddy/templates/caddy.xml.erb
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='mholt/caddy' type='service' version='0'>
    <create_default_instance enabled='true'/>
    <single_instance/>
    <dependency name='multi-user-server' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/multi-user-server'/>
    </dependency>

    <exec_method name='start'
                 type='method'
                 exec='/opt/local/bin/caddy -conf /config/caddy/Caddyfile'
                 timeout_seconds='0'>
      <method_context>
        <method_credential
          group='daemon'
          privileges='basic,!proc_session,!proc_info,!file_link_any,net_privaddr'
          user='caddy'/>
        <method_environment>
          <envvar name='NS1_API_KEY' value='<%= @ns1_api_key %>'/>
        </method_environment>
      </method_context>
    </exec_method>

    <exec_method name='stop' type='method' exec=':kill -TERM' timeout_seconds='30'/>
    <property_group name='startd' type='framework'>
      <propval name='duration' type='astring' value='child'/>
    </property_group>
    <stability value='Unstable'/>
    <template>
      <common_name>
        <loctext xml:lang='C'>Caddy HTTP/2 web server</loctext>
      </common_name>
      <documentation>
        <doc_link name='caddyserver.com'
          uri='https://caddyserver.com' />
      </documentation>
    </template>
  </service>
</service_bundle>
